# startupscript.service
[Unit]
Description=Coffee Shop App — auto-update and start
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/coffee-shop
ExecStart=/bin/bash /opt/coffee-shop/autoupdate.sh
User=pi
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
